# supervisord.conf
[supervisord]

nodaemon=true

[program:jenkins]

command=java -jar /opt/jenkins.war
